# deploy/systemd/scholacantorum.service
[Unit]
Description=Schola Cantorum Astro app
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/var/www/scholacantorum
Environment=NODE_ENV=production
Environment=HOST=127.0.0.1
Environment=PORT=4322
ExecStart=/usr/bin/node dist/server/entry.mjs
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
